# production/magento/config/crontab/www
# run magento cronjob
*/1 * * * * ! test -e /var/www/src/var/.maintenance.flag && /usr/local/bin/php /var/www/src/bin/magento cron:run >> /var/www/src/var/log/magento.cron.log
